# configs/osquery/agentsentinel.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem,syslog",
    "logger_path": "/var/log/osquery",
    "disable_logging": "false",
    "schedule_splay_percent": "10",
    "pidfile": "/var/osquery/osquery.pidfile",
    "events_expiry": "3600",
    "database_path": "/var/osquery/osquery.db",
    "verbose": "false",
    "worker_threads": "2",
    "enable_monitor": "true"
  },
  
  "schedule": {
    "agent_processes": {
      "query": "SELECT p.pid, p.name, p.path, p.cmdline, p.cwd, p.uid, u.username, p.start_time FROM processes p LEFT JOIN users u ON p.uid = u.uid WHERE p.name LIKE '%python%' OR p.name LIKE '%node%' OR p.cmdline LIKE '%agent%'",
      "interval": 60,
      "description": "Monitor agent-related processes"
    },
    
    "sensitive_file_access": {
      "query": "SELECT * FROM file_events WHERE target_path LIKE '%/.env%' OR target_path LIKE '%/secrets%' OR target_path LIKE '%wallet%' OR target_path LIKE '%private%' OR target_path LIKE '%key.json%'",
      "interval": 30,
      "description": "Monitor access to sensitive files"
    },
    
    "env_file_integrity": {
      "query": "SELECT path, sha256, mtime FROM hash WHERE path IN ('/app/.env', '/home/agent/.env', '/etc/agentsentinel/config.yaml')",
      "interval": 300,
      "description": "Check integrity of configuration files"
    },
    
    "network_connections": {
      "query": "SELECT DISTINCT s.pid, p.name, s.remote_address, s.remote_port, s.local_port, s.state FROM process_open_sockets s JOIN processes p ON s.pid = p.pid WHERE s.remote_address NOT IN ('127.0.0.1', '::1', '0.0.0.0') AND s.remote_address != ''",
      "interval": 60,
      "description": "Monitor outbound network connections"
    },
    
    "listening_ports": {
      "query": "SELECT l.pid, p.name, l.port, l.protocol, l.address FROM listening_ports l JOIN processes p ON l.pid = p.pid",
      "interval": 300,
      "description": "Monitor listening ports"
    },
    
    "cron_jobs": {
      "query": "SELECT * FROM crontab",
      "interval": 3600,
      "description": "Monitor scheduled tasks"
    },
    
    "shell_history": {
      "query": "SELECT * FROM shell_history WHERE command LIKE '%curl%' OR command LIKE '%wget%' OR command LIKE '%nc %' OR command LIKE '%base64%'",
      "interval": 300,
      "description": "Monitor suspicious shell commands"
    },
    
    "sudo_usage": {
      "query": "SELECT * FROM last WHERE type = 7",
      "interval": 60,
      "description": "Monitor privilege escalation"
    },
    
    "docker_containers": {
      "query": "SELECT id, name, image, state, started_at FROM docker_containers",
      "interval": 120,
      "description": "Monitor Docker containers"
    },
    
    "memory_map_suspicious": {
      "query": "SELECT pid, path FROM process_memory_map WHERE path LIKE '/dev/shm%' OR path LIKE '/tmp%'",
      "interval": 300,
      "description": "Detect in-memory execution"
    }
  },
  
  "file_paths": {
    "agent_configs": [
      "/app/.env",
      "/app/config/%%",
      "/home/agent/.env"
    ],
    "wallet_files": [
      "/home/agent/.config/solana/%%",
      "/app/wallets/%%"
    ],
    "ssh_keys": [
      "/home/%%/.ssh/%%"
    ]
  },
  
  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT hostname FROM system_info;"
    ]
  }
}
